# conbee.service
[Unit]
Description=Conbee II Prometheus exporter

[Service]
ExecStart=/usr/local/bin/conbee-exporter
